# base/etc/update-motd.d/10-uname
#!/bin/sh
echo "\033[1;36m	Mier Engineering \033[1;34mSystem \033[1;32mU\033[1;34mX "
cat \etc\UX_Version
echo " \033[1;32m"
echo ""
echo "Kernel:"
uname -rm
echo " \033[1;0m"
