# SRC/base/etc/update-motd.d/10-uname
#!/bin/sh
echo -e "\033[1;34m	Mier Engineering NT-UX V.0.01.08 281222 \033[1;32m"
echo ""
echo "Kernel:"
uname -rvm
echo -e " \033[1;0m"
