# SRC/base/etc/update-motd.d/10-uname
#!/bin/sh
echo -e "\033[1;34m	Mier Engineering System \033[1;32mU\033[1;34mX V.0.01.09 281222 \033[1;32m"
echo ""
echo "Kernel:"
uname -rvm
echo -e " \033[1;0m"
